# .launchd/com.powerdemand.autocommit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.powerdemand.autocommit</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ibkim/Ormi_1/power-demand-forecast/.git/hooks/auto-commit.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ibkim/Ormi_1/power-demand-forecast</string>

    <key>StartInterval</key>
    <integer>1800</integer> <!-- 30분마다 실행 -->

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/ibkim/Ormi_1/power-demand-forecast/logs/autocommit.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ibkim/Ormi_1/power-demand-forecast/logs/autocommit.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
